# save-attachments.crontab
*/5 * * * * mail /opt/save-attachments.sh >> /var/mail/save-attachments.log 2>&1

